# configs/supervisor/supervisor.d/php7.ini
[program:php-fpm]
command=/usr/sbin/php-fpm7 -c /etc/php7/php-fpm.conf --nodaemonize
autostart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
